# autoobs/obsheets.service
[Unit]
Description=Observation Sheets (Bridgwater ASC)
After=network.target

[Service]
# Edit WorkingDirectory/ExecStart paths to where you put the autoobs folder.
WorkingDirectory=/opt/autoobs
ExecStart=/opt/autoobs/venv/bin/gunicorn -w 1 --threads 8 -b 0.0.0.0:8000 app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
